# systemd/netevd.service
# SPDX-License-Identifier: LGPL-3.0-or-later

[Unit]
Description=Network Event Daemon
Documentation=man:netevd(8)
After=network-pre.target systemd-networkd.service
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/bin/netevd
Restart=on-failure
RestartSec=5s

# Security settings
User=netevd
Group=netevd

# Capabilities needed for network configuration
# Only CAP_NET_ADMIN is required for routing tables and netlink operations
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/netevd
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
